# OS/ruflo-os/init-system/service-manager/ruflo-agent.service
[Unit]
Description=Ruflo Agent Runtime
After=nemoclaw.service
Requires=nemoclaw.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /app/ruflo-agent/core/agent_runtime.py
Restart=on-failure
RestartSec=5s
User=ruflo
Group=ruflo
Environment=DISPLAY=:0
Environment=WAYLAND_DISPLAY=wayland-0
Environment=XDG_RUNTIME_DIR=/run/user/1000

[Install]
WantedBy=multi-user.target