# distro/systemd/ruflo-control-plane.service
[Unit]
Description=Ruflo Control Plane — Task orchestration and policy enforcement
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=exec
User=ruflo
Group=ruflo
WorkingDirectory=/opt/ruflo/control-plane
ExecStart=/opt/ruflo/venv/bin/uvicorn ruflo_control_plane.main:app \
    --host 127.0.0.1 \
    --port 9000 \
    --workers 2
EnvironmentFile=/etc/ruflo/control-plane.env
Restart=on-failure
RestartSec=5
TimeoutStartSec=30

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/ruflo
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

[Install]
WantedBy=multi-user.target
